# scripts/com.jarvis.dashboard.template.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jarvis.dashboard</string>

    <key>WorkingDirectory</key>
    <string>__JARVIS_HOME__</string>

    <key>ProgramArguments</key>
    <array>
        <string>__NPM_PATH__</string>
        <string>start</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__NODE_DIR__:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>ProcessType</key>
    <string>Interactive</string>

    <key>StandardOutPath</key>
    <string>__USER_HOME__/Library/Logs/jarvis.out.log</string>

    <key>StandardErrorPath</key>
    <string>__USER_HOME__/Library/Logs/jarvis.err.log</string>
</dict>
</plist>
